# dat/baselayout/etc/inittab
::sysinit:/etc/init.d/rcS start
hvc0::askfirst:/bin/sh
::shutdown:/etc/init.d/rcS stop
